# wmbusmeters-ha-addon-edge/rootfs/etc/services.d/wmbusmeters/finish
#!/usr/bin/env bashio
# ==============================================================================
# Take down the S6 supervision tree when service fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

PROCESS_NAME="wmbusmeters"

kill_process() {
    if pgrep "$PROCESS_NAME" > /dev/null; then
        pkill "$PROCESS_NAME"       
        sleep 1
        
        if pgrep "$PROCESS_NAME" > /dev/null; then
            pkill -9 "$PROCESS_NAME"
        fi
    fi
}

# Check exit status and take action accordingly
if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
    sleep 1
    kill_process
    exit 0
fi
exit 0
