# wmbusmeters-ha-addon/rootfs/etc/services.d/wmbusmeters/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start wmbusmeters service
# ==============================================================================

bashio::log.info "Starting wmbusmeters service."
exec /run.sh
