# wmbusmeters-ha-addon/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start NGINX service
# ==============================================================================

# Wait for backend to become available
bashio::net.wait_for 5000 localhost 900

bashio::log.info "Starting Nginx..."
exec nginx
