# wmbusmeters-ha-addon-edge/rootfs/etc/services.d/availability/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start MQTT availability monitor
# ==============================================================================

bashio::log.info "Starting availability service."
exec /availability.sh
